# infra/launchd/com.brickbot.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brickbot.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jonbrick/projects/brickbot/infra/launchd/brickbot-daily.sh</string>
    </array>

    <!-- Run at 8:00 AM daily -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If Mac was asleep at 8am, run when it wakes up -->
    <key>StartInterval</key>
    <integer>0</integer>

    <key>StandardOutPath</key>
    <string>/Users/jonbrick/projects/brickbot/local/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jonbrick/projects/brickbot/local/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/jonbrick/projects/brickbot</string>
</dict>
</plist>
