# infra/launchd/com.brickbot.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brickbot.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jonbrick/projects/brickbot/infra/launchd/brickbot-daily.sh</string>
    </array>

    <!-- Run 3x daily: 8am, 10am, 7pm -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>19</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/jonbrick/projects/brickbot/local/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jonbrick/projects/brickbot/local/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/jonbrick/projects/brickbot</string>
</dict>
</plist>
