# infra/launchd/com.brickbot.daily.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brickbot.daily</string>

    <!-- Wraps `node cli/sync.js --auto` in caffeinate so the system holds awake -->
    <!-- for exactly the script's lifetime. Wakelock releases on script exit.   -->
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/projects/brickbot/scripts/run-with-wakelock.sh</string>
        <string>/opt/homebrew/bin/node</string>
        <string>__HOME__/projects/brickbot/cli/sync.js</string>
        <string>--auto</string>
    </array>

    <!-- 9 runs/day — every 2 hours, 7 AM–11 PM (wake to sleep).                   -->
    <!-- Paired with pmset wake schedule (5 min before each).                       -->
    <!--  7 AM — kickoff (concurrent with 7 AM Cowork morning-brief)               -->
    <!--  9, 11 AM, 1, 3, 5, 7 PM — daytime cadence                                -->
    <!--  9 PM — concurrent with 9 PM Cowork evening-processor                     -->
    <!--  11 PM — end of day                                                       -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>__HOME__/projects/brickbot/local/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/projects/brickbot/local/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>__HOME__/projects/brickbot</string>
</dict>
</plist>
