# infra/launchd/com.brickbot.pull-linear.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brickbot.pull-linear</string>

    <!-- Wraps `node cli/pull-linear.js` in caffeinate so the system holds awake -->
    <!-- for exactly the script's lifetime. Wakelock releases on script exit.    -->
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/projects/brickbot/scripts/run-with-wakelock.sh</string>
        <string>/opt/homebrew/bin/node</string>
        <string>__HOME__/projects/brickbot/cli/pull-linear.js</string>
    </array>

    <!-- Single fire at 6:30 AM — deliberately separate from `yarn sync` so a  -->
    <!-- Linear failure can't stale the Notion/Calendar caches. Refreshes      -->
    <!-- data/workProjects.json + data/workTasks.json for the day's sessions.  -->
    <!-- Clear of daily-brief (6:00) and the 7:00 sync/Cowork window.          -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>6</integer>
        <key>Minute</key><integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/projects/brickbot/local/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/projects/brickbot/local/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>__HOME__/projects/brickbot</string>
</dict>
</plist>
